# supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/copy-perp-supervisord.log
pidfile=/tmp/copy-perp-supervisord.pid

[unix_http_server]
file=/tmp/copy-perp-supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/copy-perp-supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:copy-perp]
command=python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8001 --log-level warning --workers 1
directory=/root/.openclaw/workspace/paperclip-company/projects/pacifica-hackathon/copy-perp
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=30
stdout_logfile=/tmp/copy-perp.log
stderr_logfile=/tmp/copy-perp-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1",NETWORK="mainnet",ACCOUNT_ADDRESS="3AHZqrocSguMuo9sUUP8G8YN8NwHwWV2DPUQvbDvtfaQ",AGENT_WALLET="9mxJJAQwKLmM3hUdFebFXgkD8TPnDEJCZWhWN2uLZHWi",BUILDER_CODE="noivan",BUILDER_FEE_RATE="0.001",PRIVY_APP_ID="cmmvoxcix058e0ckv7uhp9ip0",PRIVY_APP_SECRET="privy_app_secret_38XQZ7nZ8CMEuHEqsVADisWqdtuQTMcLmug33DWUWggCLzUXfAz1swaEavGvmM2haC1rZoa83UPfiWbmAEESXixG",PRIVY_JWKS_URL="https://auth.privy.io/api/v1/apps/cmmvoxcix058e0ckv7uhp9ip0/jwks.json",AGENT_PRIVATE_KEY="3RtDRB8ZWVdcKfmqa36zJVyKjEL4m75RPg5weLTpZKLBhpc5yHoAzJrQSUUDGnwUJZHhjspVL3jCyh76aUN7P3uk",FUUL_API_KEY=""

[program:mainnet-tracker]
command=python3 scripts/mainnet_continuous_collector.py --loop
directory=/root/.openclaw/workspace/paperclip-company/projects/pacifica-hackathon/copy-perp
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stdout_logfile=/tmp/mainnet-tracker.log
stderr_logfile=/tmp/mainnet-tracker-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1",NETWORK="mainnet"
