# supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/copy-perp-supervisord.log
pidfile=/tmp/copy-perp-supervisord.pid

[unix_http_server]
file=/tmp/copy-perp-supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/copy-perp-supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:copy-perp]
command=python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8001 --log-level warning --workers 1
directory=/root/.openclaw/workspace/paperclip-company/projects/pacifica-hackathon/copy-perp
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=30
stdout_logfile=/tmp/copy-perp.log
stderr_logfile=/tmp/copy-perp-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB
; ⚠️ SECURITY: 시크릿은 절대 여기에 하드코딩 금지.
; 환경변수는 Render/Railway 대시보드 또는 .env 파일로 주입하세요.
; supervisord는 .env를 자동으로 읽지 않으므로, start.sh 등에서
; export $(grep -v '^#' .env | xargs) 후 supervisord를 실행하거나
; 배포 플랫폼의 환경변수 설정 UI를 사용하세요.
environment=PYTHONUNBUFFERED="1"

[program:mainnet-tracker]
command=python3 scripts/mainnet_continuous_collector.py --loop
directory=/root/.openclaw/workspace/paperclip-company/projects/pacifica-hackathon/copy-perp
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=30
stdout_logfile=/tmp/mainnet-tracker.log
stderr_logfile=/tmp/mainnet-tracker-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1",NETWORK="mainnet"

[program:pt-engine]
command=python3 papertrading/pt_engine.py
directory=/root/.openclaw/workspace/paperclip-company/projects/pacifica-hackathon/copy-perp
autostart=true
autorestart=true
startretries=10
startsecs=10
stopwaitsecs=30
stdout_logfile=/tmp/pt_engine.log
stderr_logfile=/tmp/pt_engine_err.log
stdout_logfile_maxbytes=30MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1"
