# supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/copy-perp-supervisord.log
pidfile=/tmp/copy-perp-supervisord.pid

[unix_http_server]
file=/tmp/copy-perp-supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/copy-perp-supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:copy-perp]
command=python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8001 --log-level warning --workers 1
directory=/root/.openclaw/workspace/paperclip-company/projects/pacifica-hackathon/copy-perp
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=30
stdout_logfile=/tmp/copy-perp.log
stderr_logfile=/tmp/copy-perp-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1",NETWORK="testnet",ACCOUNT_ADDRESS="3AHZqrocSguMuo9sUUP8G8YN8NwHwWV2DPUQvbDvtfaQ",AGENT_WALLET="9mxJJAQwKLmM3hUdFebFXgkD8TPnDEJCZWhWN2uLZHWi",BUILDER_CODE="noivan",BUILDER_FEE_RATE="0.001"
